# build-t2503180/apps/file/etc/supervisor.d/virtlogd.conf
[program:virtlogd]
command=/usr/sbin/virtlogd
priority=201
autostart=true
autorestart=true
stdout_logfile=/var/log/virtlogd_out.log
stderr_logfile=/var/log/virtlogd_err.log
